# mongo/supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/opt/mongodb-linux-x86_64-2.6.5/bin/mongod --config /opt/mongodb-linux-x86_64-2.6.5/mongo.conf
