# redis/supervisor.conf
[supervisord]
nodaemon=true

[program:redis]
command=/opt/redis-2.8.17/src/redis-server /opt/redis-2.8.17/redis.conf
